# deploy/systemd/equity-options-daily.service
[Unit]
Description=Daily Equity Options container startup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/equity-options
ExecStart=/opt/equity-options/scripts/daily-startup.sh
TimeoutStartSec=6h
